# deploy/systemd/neuroedge-core.target
[Unit]
Description=NeuroEdge Core Services Target
Wants=neuroedge-kernel.service neuroedge-ml.service neuroedge-orchestrator.service
After=neuroedge-kernel.service neuroedge-ml.service neuroedge-orchestrator.service
